# dockerlab.conf
# dockerlab.conf
start on filesystem
script
    cd /opt/dockerlab
    exec /opt/dockerlab/DockerLab.py
end script

